# script/com.kittentts.wyoming.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.kittentts.wyoming</string>
    <key>ProgramArguments</key>
    <array>
        <string><SCRIPT-DIR>/service</string>
        <string><SHIM></string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string><LOG-DIR>/stdout.log</string>
    <key>StandardErrorPath</key>
    <string><LOG-DIR>/stderr.log</string>
</dict>
</plist>
